# autoupdates.service
[Unit]
Description=Eopkg Up

[Service]
ExecStartPre=/usr/bin/eopkg ur
ExecStart=/usr/bin/eopkg up -y
